# crontab
MAILTO="greg.nicholas@gmail.com"
0 0 * * *	PYTHONPATH=/home/dotcloud/current SECRET_KEY=abcdefgh /home/dotcloud/env/bin/python /home/dotcloud/current/mwcdc/scripts/syncplaces.py
